# fetch.conf
description "Fetch Application"
author "James Powell"
start on runlevel [2345]
stop on runlevel [!2345]
respawn
script
cd /home/Fetch
sudo nohup java -jar *.jar &
end script